# besnap.cron
#!/bin/sh

# for cron.daily
exec /usr/local/sbin/besnap
